# supervisor_conf.d/apache2.conf
[program:apache2]
command=apachectl -e info -DFOREGROUND -e info
priority=997
user=root
autostart=true
autorestart=true
startretries=1
